# asus-kbd-backlight.conf
# Ubuntu upstart file

respawn
respawn limit 15 5

start on runlevel [2345]
stop on runlevel [06]

script
  /usr/bin/asus-kbd-backlight allowusers
end script
